# root/etc/s6-overlay/s6-rc.d/init-dnscrypt-config/up
/etc/s6-overlay/scripts/init-dnscrypt-config
